# .config/restic/restic-backup.timer
[Unit]
Description=Local Restic backup service

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target